# podversion.socket
[Unit]
Description=Socket for podversion

[Socket]
ListenStream=[::]:3080

[Install]
WantedBy=sockets.target
