# podversion.socket
[Unit]
Description=Socket for podversion

[Socket]
ListenStream=/run/podversion.sock

[Install]
WantedBy=sockets.target
